# debbuild/debian/prerm
#!/bin/sh

systemctl disable gothing.service >/dev/null
systemctl stop gothing.service >/dev/null