# debbuild/debian/prerm
#!/bin/sh

deb-systemd-helper disable gothing.service >/dev/null
deb-systemd-invoke stop gothing.service >/dev/null